# elastic-stack/backup/elasticsearch_backup.cron
@midnight   root    DATE=`date +\%Y-\%m-\%d` && /usr/bin/nodejs /var/lib/elasticsearch-backup/backup.js >> /var/log/elasticsearch-backup/${DATE}.log 2>&1
