# elastic-stack/amazon-backup/elastic_db_backup.cron
@midnight   root    DATE=`date +\%Y-\%m-\%d` && /usr/bin/nodejs /var/lib/elasticsearch-backup/elastic_db_backup.js >> /var/log/elasticsearch-backup/${DATE}.log 2>&1
